# scripts/install/vga2usbPermissions/85-paol_rule.rules
ACTION=="add", \
SUBSYSTEM=="usb", \
ATTRS{idVendor}=="5555", \
RUN+="/home/paol/paol-code/scripts/install/vga2usbPermissions/chmodBash"